# wm2026.service
[Unit]
Description=WM 2026 Tippspiel
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/wm2026
Environment=PATH=/home/YOUR_USER/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=SECRET_KEY=change-me-to-a-random-secret
ExecStart=/home/YOUR_USER/.local/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
